# systemd-services/services/user-cache/source/usr-cache.mount
[Unit]
Description=Mount Unit for User Cache

[Mount]
What=tmpfs
Where=/usr/cache
Type=tmpfs
Options=noswap

[Install]
WantedBy=multi-user.target
